# scripts/com.oliver.testschule.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oliver.testschule</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__PROJECT_DIR__/node_modules/.bin/next</string>
        <string>start</string>
        <string>--port</string>
        <string>3000</string>
        <string>--hostname</string>
        <string>0.0.0.0</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/data/server.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/data/server.err.log</string>
</dict>
</plist>
